# stuff/ftn
#! /bin/sh
#
# Start FTN services.
#
# This is a script for RedHat's startup process.
#
# $Id: ftn,v 1.1.1.1 2003/07/12 21:27:19 sisoft Exp $
# chkconfig: 345 86 14
# description: FTN services
#
. /etc/rc.d/init.d/functions

case "$1" in
  start)
        echo $PATH
	echo -n "Starting FTN services: "
	export LANG=ru # dirty hack for russian in log 
	su news -c 'qico -d'
	echo qico
	touch /var/lock/subsys/fido
	;;
  stop)
	echo -n "Stopping FTN services: "
	qctl -q
	echo qico
	rm -f /var/lock/subsys/fido
	;;
  status)
	status qico
	;;
  restart|reload)
	qctl -R
	;;
  *)
	echo "Usage: ftn {start|stop|status|restart|reload}"
	exit 1
esac

exit 0
